# templates/pkcs11_eventmgr.service.erb
[Unit]
Description=SmartCard PKCS#11 Event Manager
Documentation=man:pkcs11_eventmgr(1)

[Service]
ExecStart=/usr/bin/pkcs11_eventmgr

[Install]
WantedBy=default.target
